# Coding_Experiments/Deployment/com.algo_trading.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.alphajet_io.plist</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/fdautil</string>
		<string>exec</string>
		<string>/usr/local/Caskroom/miniforge/base/envs/dev310/bin/python</string>
		<string>/Users/adriencaudron/CODE/adrienCAD.github.io/upload_console_output.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/adriencaudron/CODE/adrienCAD.github.io/upload_console_output_error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/adriencaudron/CODE/adrienCAD.github.io/upload_console_output.log</string>
	<key>StartInterval</key>
	<integer>3600</integer>
</dict>
</plist>
